# G25_root/etc/udev/rules.d/80-door_access.rules
ACTION=="add" SUBSYSTEMS=="usb", ENV{ID_MODEL}=="*Flash_Disk*" RUN+="/usr/local/door_access/update_permissions.sh"
