# scripts/80-door_access.rules
ACTION=="add" SUBSYSTEMS=="usb", ENV{ID_MODEL}=="*Flash_Disk*" RUN+="/home/acme/door_access/scripts/update_permissions.sh"
